# deploy/lessenza.service
; systemd unit for L'Essenza Beauty Salon self-hosted server
; Install:  sudo cp deploy/lessenza.service /etc/systemd/system/
;           sudo systemctl daemon-reload
;           sudo systemctl enable --now lessenza
; Logs:     journalctl -u lessenza -f

[Unit]
Description=L'Essenza Beauty Salon (Node + Express + SQLite)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lessenza
Group=lessenza
WorkingDirectory=/opt/lessenza/app
EnvironmentFile=/opt/lessenza/app/.env
ExecStart=/usr/bin/env node --no-warnings --import tsx/esm server/index.ts
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

; Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/opt/lessenza/app/data

; Resource guardrails — salon site is tiny, but cap anyway.
LimitNOFILE=4096
MemoryMax=512M

[Install]
WantedBy=multi-user.target
